# bash/.config/profile.d/00-helper-functions
in_path() {
	hash "$1" 2>/dev/null
}

# vim: ft=sh
