# bash/.config/profile.d/10-local-path
# Add the local bin to the PATH
[ -d $HOME/.local/bin ] && PATH="$HOME/.local/bin:$PATH"

# vim: ft=sh
