# bash/.config/profile.d/30-desktops
# Export XDG paths for legacy desktop environments
export KDEHOME="$XDG_DATA_HOME/kde4"

# vim: ft=sh
