# bash/.config/profile.d/80-tex
# Export XDG paths for TeX applications
export TEXMFCONFIG="$XDG_CONFIG_HOME/texmf"
export TEXMFHOME="$XDG_DATA_HOME/texmf"
export TEXMFVAR="$XDG_CACHE_HOME/texmf"

# vim: ft=sh
